# games-util/arrpc/files/arrpc.service
[Unit]
Description=Open Discord RPC server for atypical setups

[Service]
ExecStart=/usr/bin/arrpc

[Install]
WantedBy=default.target
